# crypto-monitor-backend.service
[Unit]
Description=Crypto Monitor Backend Service
After=network.target

[Service]
Type=simple
User=dango
WorkingDirectory=/home/dango/CryptoMonitoringbot
Environment=PYTHONPATH=/home/dango/CryptoMonitoringbot
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/dango/CryptoMonitoringbot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=VIRTUAL_ENV=/home/dango/CryptoMonitoringbot/venv
Environment=OPENROUTER_API_KEY=your_openrouter_api_key
ExecStart=/home/dango/CryptoMonitoringbot/venv/bin/python /home/dango/CryptoMonitoringbot/run_backend.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 
